# platforms/alpine/templates/periodic/daily/rsnapshot
#!/bin/bash

exec /usr/bin/rsnapshot beta
