# platforms/alpine/templates/periodic/weekly/rsnapshot
#!/bin/bash

exec /usr/bin/rsnapshot gamma
