# platforms/alpine/templates/periodic/monthly/rsnapshot
#!/bin/bash

exec /usr/bin/rsnapshot delta
