# platforms/alpine/templates/periodic/hourly/rsnapshot
#!/bin/bash

exec /usr/bin/rsnapshot alpha
