# alpine-linux/templates/init/container_komodo
#!/sbin/openrc-run

name="komodo"
description="container - komodo"

restart_command="/usr/bin/docker compose -f /opt/core-services/komodo/compose.yaml restart"
start_command="/usr/bin/docker compose -f /opt/core-services/komodo/compose.yaml up -d"
stop_command="/usr/bin/docker compose -f /opt/core-services/komodo/compose.yaml down"
stop_command_args="--rmi all --volumes"

depend() {
        need container_id
        use docker
}

start() {
        ebegin "STARTING - $name"
        ${start_command}
        eend $?
}

stop() {
        ebegin "STOPPING - $name"
        ${stop_command} ${stop_command_args}
        eend $?
}

restart() {
        ebegin "RESTARTING - $name"
        ${restart_command}
        eend $?
}
