# alpine-linux/templates/periodic/weekly/rsnapshot
#!/bin/sh

exec /usr/bin/rsnapshot gamma
